# systemd/vnc@.service
[Unit]
Description=VNC Desktop :%i
After=network.target

[Service]
Type=forking
User=aadm
Group=aadm
Environment=HOME=/var/lib/aadm
ExecStart=/usr/bin/vncserver :%i -localhost -SecurityTypes None
ExecStop=/usr/bin/vncserver -kill :%i
Restart=on-failure

[Install]
WantedBy=multi-user.target
